# edit_systemd/docker/flag.service
[Unit]
Description=Flag Verify Service

[Service]
Type=oneshot
User=root
ExecStart=/usr/local/lib/flag.sh

[Install]
WantedBy=multi-user.target